# projects/bitcoin_trader_llm/systemd/start_paper.timer
[Unit]
Description=Timer to start 2-week paper trading

[Timer]
OnCalendar=*-*-* *:00:00
# OnCalendar can be customized by the operator. This is a placeholder.
Persistent=true

[Install]
WantedBy=timers.target
