# projects/bitcoin_trader_llm/systemd/start_paper.service
[Unit]
Description=Start 2-week paper trading (placeholder)
After=network.target

[Service]
Type=oneshot
ExecStart=/root/.openclaw/workspace/projects/bitcoin_trader_llm/scripts/start_paper.sh
WorkingDirectory=/root/.openclaw/workspace/projects/bitcoin_trader_llm
User=$USER
